# deploy/systemd/otedama.service
[Unit]
Description=Otedama Mining Pool
Documentation=https://github.com/otedama/otedama
After=network.target bitcoind.service
Wants=bitcoind.service

[Service]
Type=simple
ExecStartPre=/usr/bin/npm install --production
ExecStart=/usr/bin/node /opt/otedama/start-mining-pool-enhanced.js --mode hybrid
ExecReload=/bin/kill -USR1 $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# Process management
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# User and permissions
User=otedama
Group=otedama
UMask=0027

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/otedama/data /opt/otedama/logs /opt/otedama/backups
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateMounts=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=4G
CPUQuota=200%

# Environment
Environment="NODE_ENV=production"
Environment="NODE_OPTIONS=--max-old-space-size=4096"
EnvironmentFile=/opt/otedama/.env
WorkingDirectory=/opt/otedama

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=otedama

[Install]
WantedBy=multi-user.target
